# .docker/php/profile.d/aliases.sh
alias cc='php bin/console cache:clear'
alias xphp='PHP_IDE_CONFIG="serverName=develop" php -dxdebug.start_with_request=yes'
